# install/desktop/etc/profile.d/qt-logging.sh
QT_LOGGING_RULES='*.debug=false'
export QT_LOGGING_RULES
